# Lesson 5/vlad.kharlamov/wordpress/supervisord.conf
[supervisord]
nodaemon=true

[program:php7]
command=php-fpm7.2 -R -F
autostart=true
autorestart=true

[program:nginx]
command=nginx  -g 'daemon off;'
autostart=true
autorestart=true
